# profile.d/chrootprompt

# set variable identifying the chroot you work in (used in the prompt below)
if [ -z "$debian_chroot" ] && [ -r /etc/debian_chroot ]; then
    debian_chroot=$(cat /etc/debian_chroot)
fi

#export PS1="\[\e]0;\h: \w\a\]${debian_chroot:+($debian_chroot)}\[\033[01;32m\]\h\[\033[00m\]:\w"
export PS1="$PS1\[\033[01;32m\]${debian_chroot:+($debian_chroot)}\[\033[00m\]"

