# profile.d/debhelper
export DEBFULLNAME="Mark Mims"
export DEBEMAIL="mark.mims@gmail.com"
