# profile.d/keys

#[ -z $DISPLAY ] || xmodmap ~/.Xmodmap

#pgrep xcape || /usr/bin/xcape -e "Alt_L=Tab"
