# profile.d/truncate-path-in-prompt


# shorter path length in the prompt
export PROMPT_DIRTRIM=3

